# data/systemd/udisks.service.in
[Unit]
Description=Disk Manager (legacy version)
Documentation=man:udisks(7)

[Service]
Type=dbus
BusName=org.freedesktop.UDisks
ExecStart=@libexecdir@/udisks-daemon --no-debug
